# roles/monitoring/templates/tuxedovpn-cert-expiry.service.j2
[Unit]
Description=TuxedoVPN TLS cert expiry metrics (node_exporter textfile)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tuxedovpn_cert_expiry_metrics.sh
